# packaging/xui-pro-agent.service
[Unit]
Description=XUI Pro Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/xui-pro/agent.env
ExecStart=/usr/local/xui-pro/xui-agent
Restart=on-failure
RestartSec=5s
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
