# packaging/xui-pro.service
[Unit]
Description=XUI Pro Master
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/xui-pro/master.env
ExecStart=/usr/local/xui-pro/xuid
Restart=on-failure
RestartSec=5s
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
